# server/package/postinst
systemctl daemon-reload

# Make the service start at boot
update-rc.d rpis defaults